# systemd/life-agent.service.template
[Unit]
Description=Life Agent — periodic webcam snapshot + LLM reasoning
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{BUN_PATH}} run src/index.ts --config {{PROJECT_DIR}}/config.local.yml

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=life-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths={{PROJECT_DIR}}/logs {{PROJECT_DIR}}/captures {{PROJECT_DIR}}/memory
PrivateTmp=true

[Install]
WantedBy=default.target
